# .devcontainer/s6/nix-daemon/run
#!/command/execlineb
fdmove -c 2 1
nix-daemon
